# auto-reply-bot.service
[Unit]
Description=Auto Reply Bot for Naver Cafe
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/auto_reply
Environment=HEADLESS=true
Environment=PYTHONUNBUFFERED=1

# 시작 전 Chrome 좀비 프로세스 정리
ExecStartPre=/bin/bash -c 'pkill -9 chrome 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'pkill -9 chromedriver 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'rm -rf /home/azureuser/auto_reply/chrome_data_* 2>/dev/null || true'
ExecStartPre=/bin/sleep 2

# 봇 실행
ExecStart=/home/azureuser/auto_reply/run_bot.sh

# 종료 시 정리
ExecStopPost=/bin/bash -c 'pkill -9 chrome 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'rm -rf /home/azureuser/auto_reply/chrome_data_* 2>/dev/null || true'

# 재시작 설정
Restart=on-failure
RestartSec=30

# 로그
StandardOutput=append:/home/azureuser/auto_reply/bot.log
StandardError=append:/home/azureuser/auto_reply/bot.log

[Install]
WantedBy=multi-user.target
